# Buda/Docker/Cron/cron.d
0 1 * * * python /project/Buda/Buda/cron_buda.py >> /var/log/buda.log