# Buda/Docker/Cron/cron.d
# min	hour	day		month	weekday		command
0		1		*		*		*			/usr/local/bin/python /project/Buda/Buda/cron_buda.py >> /var/log/buda.log
*/5		*		*		*		*			/usr/local/bin/python /project/Buda/Buda/cron_buda.py >> /var/log/buda.log